# whisper_transcription_server/whisper.service
[Unit]
Description=Whisper Transcription Server
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/user/whisper_server.py
WorkingDirectory=/home/user
Environment="OPENAI_API_KEY=your-api-key-here"
Restart=always
User=user

[Install]
WantedBy=multi-user.target
